# etc/flux-housekeeping@.service.in
[Unit]
Description=Housekeeping for Flux job %I
CollectMode=inactive-or-failed

[Service]
Type=oneshot
EnvironmentFile=-@X_RUNSTATEDIR@/flux-housekeeping@%I.env
ExecStart=@X_SYSCONFDIR@/flux/system/housekeeping
ExecStopPost=-rm -f @X_RUNSTATEDIR@/flux-housekeeping@%I.env
ExecStopPost=-sh -c 'test $SERVICE_RESULT = "success" || flux resource drain $(flux getattr rank) "housekeeping $EXIT_CODE $EXIT_STATUS"'
